# internet-watcher.timer
[Unit]
Description=Failsafe periodic wakeup

[Timer]
OnBootSec=1min
OnUnitActiveSec=3min
AccuracySec=1min

[Install]
WantedBy=timers.target
