# internet-watcher.service
[Unit]
Description=Internet Watcher (event driven)

[Service]
ExecStart=/usr/bin/python3 /usr/lib/internet-watcher/daemon.py
Restart=on-failure
